# apiserver.service
[Unit]
Description=NDS Labs API Server

[Service]
ExecStart=/home/core/apiserver/apiserver -v 4 -conf /home/core/apiserver/apiserver.conf
