# ythub_worker/etcdl.conf
[program:etcd_listener]
command = /usr/local/bin/etcd_listener
autostart = true
autorestart = true
redirect_stderr = true
