# irods-idrop2/common/supervisord.conf
[supervisord]
#nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=false
autorestart=false

[program:tomcat6]
command=/home/admin/tcstart.sh
autorestart=false
autostart=true
